# docker/cron/crontab
*/1 * * * * /usr/local/bin/php /var/www/html/cron.php
*/1 * * * * /usr/local/bin/php /updatenodes.php
